# dist/porchlight-0.3.0/systemd/porchlight-discover.timer
[Unit]
Description=Run cheap Porchlight host discovery

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
RandomizedDelaySec=90s
Unit=porchlight-discover.service
Persistent=true

[Install]
WantedBy=timers.target
