# dist/porchlight-2.3.1/systemd/porchlight-ai-analysis.timer
[Unit]
Description=Refresh Porchlight AI analysis

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Unit=porchlight-ai-analysis.service
Persistent=true
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
