# dist/porchlight-0.3.1/systemd/porchlight-ha-mqtt-bridge.timer
[Unit]
Description=Refresh Porchlight Home Assistant MQTT state

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Unit=porchlight-ha-mqtt-bridge.service
Persistent=true

[Install]
WantedBy=timers.target
