# dist/porchlight-1.1.1/systemd/porchlight-health.timer
[Unit]
Description=Refresh Porchlight health status

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=porchlight-health.service
Persistent=true

[Install]
WantedBy=timers.target
