# dist/porchlight-1.1.2/systemd/porchlight-render.timer
[Unit]
Description=Refresh Porchlight rendered dashboard snapshots

[Timer]
OnBootSec=90s
OnUnitActiveSec=5min
Unit=porchlight-render.service
Persistent=true

[Install]
WantedBy=timers.target
